# debian/overlay/etc/systemd/system/bluetooth.service.d/after-bluetooth-start.conf
# Ensure bluetoothd starts after HCI is brought up by bluetooth-start.service (rtk_hciattach)
[Unit]
After=bluetooth-start.service
Wants=bluetooth-start.service
